# deploy/enoch-notion-sync.service
[Unit]
Description=Enoch Notion intake/projection sync
After=enoch-control-plane.service network-online.target
Wants=enoch-control-plane.service network-online.target

[Service]
Type=oneshot
Environment=ENOCH_INSTALL_DIR=/opt/enoch-agentic-research-system
Environment=OMX_WAKE_GATE_CONFIG=/etc/enoch/config.json
EnvironmentFile=-/etc/enoch/notion-sync.env
WorkingDirectory=/opt/enoch-agentic-research-system
ExecStart=/opt/enoch-agentic-research-system/deploy/enoch_notion_sync.sh
