# deploy/enoch-paper-draft-next.service
[Unit]
Description=Draft next eligible Enoch paper without dispatching work
After=enoch-control-plane.service
Wants=enoch-control-plane.service

[Service]
Type=oneshot
Environment=OMX_WAKE_GATE_CONFIG=/etc/enoch/config.json
WorkingDirectory=/opt/enoch-agentic-research-system
ExecStart=/opt/enoch-agentic-research-system/deploy/enoch_paper_draft_next.sh
